# SOURCES/pkg.postinst
#!/bin/bash

# We do not START these services here because we
#   depend on the universal hook to do that for us.
/bin/systemctl enable httpd.service >/dev/null 2>&1
